# server/ansible/templates/certbot-renew.timer.j2
# {{ ansible_managed }}

[Unit]
Description=renewing Let's Encrypt certificates
Requires=certbot-renew.timer

[Timer]
Unit=certbot-renew.timer
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
